# deploy/supervisor/supervisord.d/hey.ini
[group:hey]
programs=p0, p1

[program:p0]
command=python /data/hey/app.py --port=8010
directory=/data/hey/
user=www-data
autorestart=true
redirect_stderr=true
stdout_logfile=/data/supervisor/log/hey-0.log
loglevel=info

[program:p1]
command=python /data/hey/app.py --port=8011
directory=/data/hey/
user=www-data
autorestart=true
redirect_stderr=true
stdout_logfile=/data/supervisor/log/hey-1.log
loglevel=info

[supervisord]
